# deploy/systemd/slime-music-library.service
[Unit]
Description=Refresh SlimeAudio music library database
After=network-online.target mnt-rockhouse.mount

[Service]
Type=oneshot
User=squidward
WorkingDirectory=/home/squidward/.openclaw/workspace/slime-audio
ExecStart=/home/squidward/.openclaw/workspace/slime-audio/.venv/bin/python /home/squidward/.openclaw/workspace/slime-audio/scripts/slime_music_library_service.py --tunebat-backfill-limit 12 --tunebat-max-seconds 1200 --dj-analysis-backfill-limit 6 --skip-backfill-when-live --live-grace-seconds 21600
Nice=15
CPUSchedulingPolicy=batch
CPUQuota=400%
MemoryMax=8G
IOSchedulingClass=idle
